# sample_configs/rm-watcher-scheduler.timer
[Unit]
Description=Schedule rm-watcher every 6 hours

[Timer]
OnUnitActiveSec=3h
Unit=rm-watcher.service

[Install]
WantedBy=timers.target
