# sample_configs/rm-watcher.service
[Unit]
Description=A tool that downloads hot stories from reddit

[Service]
Type=simple
ExecStart=/home/your-user-name/.cargo/bin/rm-watcher
